# install/fortran-load.timer
[Unit]
Description=fortran-load

[Timer]
OnBootSec=0s
OnCalendar=hourly

[Install]
WantedBy=multi-user.target
